# frameo_control_addon/rootfs/etc/services.d/frameo/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting the Frameo Python server..."

# Read configuration and export it for the Python script
export CONNECTION_TYPE=$(bashio::config 'connection_type')
export DEVICE_SERIAL=$(bashio::config 'device_serial')
export DEVICE_HOST=$(bashio::config 'device_host')
export DEVICE_PORT=$(bashio::config 'device_port')

# Execute the Python server
exec /usr/bin/python3 /app/app.py
